# rust/05_sandbox/screensaver/own_made_screenlocker/home_dir/.config/autostart/screen locker after 5min.desktop
[Desktop Entry]
Exec=own_made_screenlocker
Name=screen locker after 5min
Type=Application
Version=1.0
