# rust/05_sandbox/dbus/custom_dbus_polkitd/root_drive/etc/local.d/custom_dbus_polkitd.start
#!/usr/bin/env bash

export RUST_BACKTRACE=1
/etc/local.d/custom_dbus_polkitd &
